# data/nscd.conf
description "name service cache daemon"

start on local-filesystems
stop on runlevel [06]

#expect fork
pre-start script
   mkdir -p /var/run/nscd
end script

exec /usr/sbin/nscd -f /etc/nscd.conf
